# ansible/roles/gunicorn/templates/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream={{ metax_base_path }}/gunicorn/socket

[Install]
WantedBy=sockets.target
